# src/rootfs/etc/supervisor/conf.d/pulseaudio.conf
[program:pulseaudio]
command=/usr/bin/pulseaudio
user=%(ENV_USER)s
autostart=true
autorestart=true
startsecs=3
startretries=1000000
priority=2
loglevel=debug
